# files/selinux/avc-audit-report.cron
# Run hourly from 6AM to 11PM EDT on weekdays (-400 UTC)
37 0-2,10-23 * * 1-5 root /usr/local/bin/avc-audit-report
# Run twice a day on weekends
37 2,10 * * 6-7 root /usr/local/bin/avc-audit-report
